# home/Library/LaunchAgents/obsidian-vault-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>task-backup</string>
    <key>UserName</key>
    <string>k</string>
    <key>ProgramArguments</key>
    <array>
      <string>bash</string>
      <string>-c</string>
			<string>cd /Users/k/vault &amp;&amp; /usr/bin/git add -A . &amp;&amp; /usr/bin/git commit -m "`/Users/k/bin/git-status-message`" &amp;&amp; /usr/bin/git push</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/launchagentout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/launchagenterr.log</string>
</dict>
</plist>
